# agents/notes_sync/com.notes_sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.notes_sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jiayi/Developer/ClaudeCode/agents/notes_sync/run.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>14400</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/jiayi/Developer/ClaudeCode/agents/notes_sync/sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jiayi/Developer/ClaudeCode/agents/notes_sync/sync.log</string>
</dict>
</plist>
